# examples/ci.crontab
MAILTO=your.email@here

## run XRef Continuous Integration every 2 minutes
*/2	*	*	*	*	xref-ci --config=/path/to/xref.ini

## run XRef doc once per day
20  3   *   *   *   xref-doc --config=/path/to/xref.ini
