# examples/ci.crontab
MAILTO=your.email@here

## this is a cron-script to run XRef Continuous Integration every 2 minutes
*/2	*	*	*	*	xref-ci --config=/path/to/your/config.file.ini
